# backup.cron
# CRONTAB entry
# backup script 
00 22 * * * root /usr/local/bin/backup/backup.sh
