# backup.cron
# CRONTAB entry
# backup script 
00 22 * * * root /opt/backup/backup.sh
